# src/deb/control/prerm
#!/bin/bash

set -e

case "$1" in
    remove|deconfigure)
        echo "Stopping HLabMonitor service..."
        systemctl stop hlabmonitor.service 2>/dev/null || true
        systemctl disable hlabmonitor.service 2>/dev/null || true
        ;;

    upgrade|failed-upgrade)
        ;;

    *)
        echo "prerm called with unknown argument \`$1'" >&2
        exit 1
        ;;
esac

exit 0
